# period_predictor/rootfs/etc/services.d/period-predictor/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Period Predictor service"

# Determine the port: prefer configured network port, then ingress, else default
PORT="$(bashio::addon.port 3002/tcp 2>/dev/null)"
if ! bashio::var.has_value "${PORT}"; then
    PORT="$(bashio::addon.ingress_port 2>/dev/null)"
fi
if ! bashio::var.has_value "${PORT}"; then
    PORT="3002"
fi
export PORT
bashio::log.info "Listening on port ${PORT}"

exec npm --prefix /app/backend start
